# elevate.conf
start on started system-services
stop on runlevel [06]
script
    /usr/local/bin/elevate-daemon
end script
